# znc.service.in
[Unit]
Description=ZNC, an advanced IRC bouncer
After=network.target

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/znc -f
User=znc

[Install]
WantedBy=multi-user.target
